# cronjobs/supervisor/boxtz.conf
[program:boxtz]
command=/usr/bin/php artisan queue:listen --env=production
directory=/var/www/boxtz
stdout_logfile=/var/www/boxtz/app/storage/logs/queue-prod.log
redirect_stderr=true

